# etc/btube.service.example
[Unit]
Description=Simple http server providing on_publish authentication for the nginx-rtmp module

[Service]
Type=simple
ExecStart=/usr/local/bin/btube-server -c /etc/btube.conf
User=btube
